# pixie/overlay/etc/inittab
# /etc/inittab — BusyBox inittab for fos-alpine initramfs
#
# OpenRC handles system initialization via runlevels.
# Autologin gettys on the serial and VGA consoles launch the imaging agent.

::sysinit:/sbin/openrc sysinit
::wait:/sbin/openrc boot
::wait:/sbin/openrc default

# Serial console — primary in headless QEMU (-nographic / -serial mon:stdio)
# Autologin runs the imaging agent.
ttyS0::respawn:/sbin/getty -L -n -l /bin/fos-autologin 115200 ttyS0 vt100

# VGA console — autologin runs the imaging agent
tty1::respawn:/sbin/getty -L -n -l /bin/fos-autologin 115200 tty1 vt100

# Emergency fallback
::ctrlaltdel:/sbin/reboot
::shutdown:/sbin/openrc shutdown
